# rootfs/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev

ttyS0::respawn:/bin/getty -L ttyS0 115200 vt100


::shutdown:/bin/umount -a -r
::shutdown:/bin/swapoff -a

